# task_10/data/logfile.service
[Unit]
Description=Logfile service

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/logfile
ExecStart=/vagrant/data/logfile.sh $WORD $LOG