# task_10/data/logfile.timer
[Unit]
Description=Logfile timer

[Timer]
OnCalendar=*:*:0/30
Unit=logfile.service

[Install]
WantedBy=multi-user.target